# spine/portals/hozan-taher/infra/launchd/com.hozan-taher.claude-context-reindex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.hozan-taher.claude-context-reindex</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/messingtomas/Documents/Projekty/hozan-taher/scripts/cron-claude-context-reindex.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/messingtomas/.cache/claude-context-reindex.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/messingtomas/.cache/claude-context-reindex-error.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/messingtomas</string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
</dict>
</plist>
