# spine/outreach/relay/deploy/anti-trace-relay.service
[Unit]
Description=Anti-Trace Relay - Privacy-hardened communication relay
After=network-online.target tor.service
Wants=network-online.target

[Service]
Type=simple
User=antitrace
Group=antitrace
WorkingDirectory=/opt/anti-trace-relay

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/anti-trace-relay/data
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Core dump prevention
LimitCORE=0

# Environment
EnvironmentFile=/opt/anti-trace-relay/env
EnvironmentFile=-/opt/anti-trace-relay/secrets

ExecStart=/opt/anti-trace-relay/anti-trace-relay
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anti-trace-relay

[Install]
WantedBy=multi-user.target
